# files/certbot-renewal.service
[Unit]
Description=Certbot Automatic Renewal Service

[Service]
Type=oneshot
ExecStart=/usr/bin/certbot renew --quiet --agree-tos